# supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:minidlna]
command=/usr/src/minidlna/minidlna-1.1.4/minidlnad -S -R -f /etc/minidlna.conf
redirect_stderr=true
