# dist/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:minidlna]
command=/minidlna/minidlnad -S -R -f /etc/minidlna.conf
redirect_stderr=true
